# .config/systemd/user/wallpaper-rotator.service
[Unit]
Description=rotate wallpaper
PartOf=graphical-session.target

[Service]
ExecStart=%h/wallpaper-rotator.sh
Restart=on-failure